# ngrok内网穿透_openVPN/natapp.ini
[program:natapp]
command                 = /prodata/ngs/bin/natapp -authtoken=xxx -log=stdout -log-level=INFO
directory               = /prodata/ngs/bin
process_name            = %(program_name)s_%(process_num)s
numprocs                = 1
autostart               = true
autorestart             = true

stdout_logfile          = /var/log/supervisor/natapp_stdout.log
stdout_logfile_maxbytes = 10MB
stderr_logfile          = /var/log/supervisor/natapp_error.log
stderr_logfile_maxbytes = 10MB
