# ngrok内网穿透_openVPN/sunny.ini
[program:sunny]
command                 = /prodata/ngs/bin/sunny clientid bf4e8a483046bbb4 vpn -log=none -log-level=INFO
directory               = /prodata/ngs/bin
process_name            = %(program_name)s_%(process_num)s
numprocs                = 1
autostart               = true
autorestart             = true

stdout_logfile          = /var/log/supervisor/sunny_std.log
stdout_logfile_maxbytes = 10MB
stderr_logfile          = /var/log/supervisor/sunny_err.log
stderr_logfile_maxbytes = 10MB
